# crontab
*/2 * * * * python3  /usr/src/app/bot_crypto.py